# configs/robair.service.template
[Unit]
Description=RobAIR

[Service]
Type=forking
PIDFile=/run/user/<UID>/robair.pid
ExecStart=<ROBAIR_HOME>/scripts/robair start --systemd

[Install]
WantedBy=default.target
